# docker/edge/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; ── soca-engine ───────────────────────────────────────────────────────────────
[program:soca-engine]
command=uvicorn main:app --host 0.0.0.0 --port 8001
directory=/app/soca-engine
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/soca-engine.log
stderr_logfile=/var/log/supervisor/soca-engine.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB

; ── soca-dashboard ────────────────────────────────────────────────────────────
[program:soca-dashboard]
command=gunicorn dashboard.wsgi:application --bind 0.0.0.0:8080 --workers 2 --threads 2 --timeout 120
directory=/app/soca-dashboard
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/supervisor/soca-dashboard.log
stderr_logfile=/var/log/supervisor/soca-dashboard.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
